# vps/systemd/digit-hab-api.service.example
# Digit-Hab API — Daphne (ASGI) sur 127.0.0.1:3004
# Projet : https://github.com/DIGIT-HABs/digit-hab-api
# ASGI : digit_hab_crm.asgi:application  (HTTP + WebSockets /ws/*)
# Settings prod : digit_hab_crm.settings.prod
#
# Prérequis : PostgreSQL + Redis (CHANNEL_LAYERS en prod utilise Redis).
# Gunicorn/WSGI ne gère PAS les WebSockets (chat, notifications).
#
# sudo cp .../digit-hab-api.service.example /etc/systemd/system/digit-hab-api.service
# sudo systemctl daemon-reload && sudo systemctl enable --now digit-hab-api.service

[Unit]
Description=Digit-Hab API (Django + Daphne ASGI sur 3004)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/apps/digit-hab-api
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=digit_hab_crm.settings.prod
EnvironmentFile=-/opt/apps/digit-hab-api/.env
ExecStart=/opt/apps/digit-hab-api/venv/bin/daphne \
  -b 127.0.0.1 \
  -p 3004 \
  digit_hab_crm.asgi:application
Restart=on-failure
RestartSec=5
MemoryMax=1G

[Install]
WantedBy=multi-user.target
